# web/docker/prod/supervisord_was.conf
[supervisord]
nodaemon=true

[program:was]
command = gunicorn -b 0.0.0.0:80 config.wsgi                 
user = root                                                           ; User to run as
stdout_logfile = /web/logs/gunicorn_supervisor.log   ; Where to write log messages
redirect_stderr = true                                                ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8                       ; Set UTF-8 as default encoding
;stdout_logfile=/dev/stdout
;stdout_logfile_maxbytes=0
;stderr_logfile=/dev/stderr
;stderr_logfile_maxbytes=0
